# rootfs/etc/s6-overlay/s6-rc.d/svc-obsidian-sync/run
#!/command/with-contenv sh
# Long-running sync service supervised by s6.

VAULT_PATH="${VAULT_PATH:-/vault}"
cd "$VAULT_PATH"

echo "[obsidian-sync] Starting continuous sync..."
exec ob sync --continuous
